# debian/bluecherry.upstart
# bc-server - Bluecherry, LLC DVR Server
#
# bc-server manages recording of cameras to disk based on rules setup in the db

description	"Bluecherry DVR Server"
author		"Ben Collins <bcollins@bluecherry.net>"

start on runlevel [23]
stop on runlevel [016]

expect fork
respawn
umask 027

pre-start script
    mkdir -p /var/run/bluecherry
    chown bluecherry.bluecherry /var/run/bluecherry
end script

exec /usr/bin/sudo -H -u bluecherry /usr/sbin/bc-server
