# debian/bluecherry.service
[Unit]
Description=Bluecherry DVR Server
After=syslog.target network.target mysql.service

[Service]
Type=simple
User=bluecherry
Group=bluecherry
#ExecStart=/usr/bin/valgrind --tool=memcheck --leak-check=yes --leak-check=full --show-leak-kinds=all --log-file=/tmp/bluecherry-valgrind.log /usr/sbin/bc-server -s
#ExecStart=/usr/bin/rr record -o /tmp/bluecherry.rr /usr/sbin/bc-server -s
ExecStart=/usr/sbin/bc-server -s
Restart=on-failure
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target

