# debian/bluecherry.cron
# m h dom mon dow user          command
0   0  *   *   *  bluecherry    /usr/share/bluecherry/backup_db.sh
30   0  *   *   *  bluecherry    curl -k https://localhost:7001/subdomainprovidercron
